# devops-ai-assistant/10_Implementation/01_vm/ansible/training-app.service
[Unit]
Description=Training Flask Application
After=network.target

[Service]
Type=simple
User=training
WorkingDirectory=/opt/training-app
EnvironmentFile=/opt/training-app/.env
ExecStart=/opt/training-app/venv/bin/python devops-ai-assistant/10_Implementation/01_vm/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
